# config/docker/etc/services.d/puma/run
#!/usr/bin/with-contenv sh
cd /app && rake db:migrate && rake db:seed && rake assets:precompile && puma
